# root/etc/services.d/tvhProxy/run
#!/usr/bin/with-contenv bash

if [[ ! -z "${TVH_URL}" ]]; then
	echo "Starting tvhProxy."
	exec \
			s6-setuidgid abc python /tvhProxy/tvhProxy.py
fi
